# rules/99-dongle-auto-restart.rules
# 99-dongle-auto-restart.rules — Auto-restart dongle after USB hotplug
# Fires when any ttyUSB device is added, triggers the auto-restart script
# in background so udev is not blocked. The script handles debounce,
# dongle.conf mapping, and the 5-second delay internally.
ACTION=="add", SUBSYSTEM=="tty", KERNEL=="ttyUSB*", RUN+="/bin/bash -c '/opt/issabel-dashboard/scripts/auto-restart-dongle.sh %k &'"
